# lib/systemd/system/dmesgd.service

[Unit]
Description=Disable dmesg output

[Service]
Type=oneshot
ExecStart=dmesg -D

[Install]
WantedBy=multi-user.target
